# deploy/douyin-backend.service
[Unit]
Description=Douyin Hackathon - Popular Science Video Deconstructor (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/douyin-ai-hackathon
# 先 sudo install -m 600 -o root -g root /dev/stdin /etc/douyin-backend.env，
# 内容：GEMINI_API_KEY=...  和  UPLOAD_DIR=/var/lib/douyin-uploads
EnvironmentFile=/etc/douyin-backend.env
ExecStart=/opt/douyin-ai-hackathon/.venv/bin/python -m uvicorn backend.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
